# etc/systemd/system/ippbxpy.timer
[Init]
Description=Timer for ippbxpy


[Timer]
RandomizedDelaySec=3500
OnCalendar=hourly

[Install]
WantedBy=timers.target
