# etc/systemd/system/ippbxpy.service
[Unit]
Description=Config generator for Asterisk and Yealink VoIP phones


[Service]
ExecStart=/usr/local/bin/ippbx.py
